# src/board/rootfs_overlay/etc/init.d/containerd
#!/sbin/openrc-run
#Service script for containerd. Documentation for writing servic scripts 
#found at https://github.com/OpenRC/openrc/blob/master/service-script-guide.md

command="/usr/local/bin/containerd"
command_args=""
pidfile="/run/${RC_SVCNAME}.pid"
command_background=true
