# src/board/rootfs_overlay/etc/init.d/kubelet
#!/sbin/openrc-run
#Service script for kubelet. Uses the supervise-daemon to start and restart the kubelet.
#Ducumentation for using the supervise-daemon for service scripts found at 
#https://github.com/OpenRC/openrc/blob/master/supervise-daemon-guide.md

supervisor=supervise-daemon

#Load environment variables set by kubeadm, i.e. $KUBELET_KUBEADM_ARGS
if [ -e /var/lib/kubelet/kubeadm-flags.env ]; then
	. /var/lib/kubelet/kubeadm-flags.env;
fi

command="/usr/bin/kubelet"
#args passed to kubelet
command_args="${KUBELET_ARGS} ${KUBELET_KUBEADM_ARGS}"
output_log="/root/log_kubelet"
error_log="/root/errorlog_kubelet"

depend() {
	after net
	need cgroups
}
